# bin/encrypt.service
[Unit]
Description=Encryption

[Service]
ExecStart=/home/iven/bin/encrypt.sh

[Install]
WantedBy=multi-user.target
